# deploy/p4e-ops/mathin-p4e-disk-check.timer
[Unit]
Description=Check Mathin disk usage and backup freshness

[Timer]
OnBootSec=5m
OnUnitActiveSec=15m
Persistent=true

[Install]
WantedBy=timers.target
