# deploy/p4e-ops/mathin-p4e-backup.timer
[Unit]
Description=Run Mathin P4E backup every day

[Timer]
OnCalendar=*-*-* 02:15:00 Asia/Shanghai
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
